# remote/production/api.service
[Unit]
Description=Elodi API service

# Ensure the service starts after PostgreSQL and network are online
After=postgresql.service
After=network-online.target
Wants=network-online.target

# Configure start limit intervals and burst
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
# Service type and user/group configuration
Type=exec
User=elodi
Group=elodi

# Load environment variables from this file
EnvironmentFile=/etc/environment

# Set working directory and executable start command
WorkingDirectory=/home/elodi
ExecStart=/home/elodi/backend/api -port=4000 -db-dsn=${ELODI_DB_DSN} -env=production

# Restart configuration
Restart=on-failure
RestartSec=5

[Install]
# Ensure the service is started in multi-user mode
WantedBy=multi-user.target